# scripts/com.fantasyhockey.weeklyreport.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fantasyhockey.weeklyreport</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/drew/Projects/yahoo-fantasy-hockey/.venv/bin/python</string>
        <string>/Users/drew/Projects/yahoo-fantasy-hockey/send_weekly_email.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/drew/Projects/yahoo-fantasy-hockey</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>  <!-- 0 = Sunday -->
        <key>Hour</key>
        <integer>8</integer>  <!-- 8 AM -->
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/drew/Projects/yahoo-fantasy-hockey/.cache/weekly_report_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/drew/Projects/yahoo-fantasy-hockey/.cache/weekly_report_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
